# docker/supervisor/supervisord.conf
; =============================================================================
; Supervisor Configuration
; =============================================================================
; Manages PHP-FPM, Nginx, and Laravel Queue Worker

[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

; =============================================================================
; PHP-FPM
; =============================================================================
[program:php-fpm]
command=/usr/local/sbin/php-fpm -F
autostart=true
autorestart=true
priority=5
stdout_logfile=/var/log/supervisor/php-fpm.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/supervisor/php-fpm-error.log
stderr_logfile_maxbytes=10MB

; =============================================================================
; Nginx
; =============================================================================
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/supervisor/nginx.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/supervisor/nginx-error.log
stderr_logfile_maxbytes=10MB

; =============================================================================
; Laravel Queue Worker
; =============================================================================
[program:laravel-queue]
command=/usr/local/bin/php /var/www/html/artisan queue:work --sleep=3 --tries=3 --max-time=3600 --memory=128
autostart=true
autorestart=true
priority=15
numprocs=2
process_name=%(program_name)s_%(process_num)02d
stdout_logfile=/var/log/supervisor/queue.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/supervisor/queue-error.log
stderr_logfile_maxbytes=10MB
stopwaitsecs=3600

; =============================================================================
; Laravel Scheduler (runs every minute)
; =============================================================================
[program:laravel-scheduler]
command=/bin/sh -c "while true; do /usr/local/bin/php /var/www/html/artisan schedule:run --verbose --no-interaction; sleep 60; done"
autostart=true
autorestart=true
priority=20
stdout_logfile=/var/log/supervisor/scheduler.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/supervisor/scheduler-error.log
stderr_logfile_maxbytes=10MB
